# scripts/ika-twitch.service
[Unit]
Description=IkaEXE Twitch Bot
After=network.target

[Service]
Type=simple
User=snaxine
WorkingDirectory=/mnt/wailord/stash/Ika's Workspace/ika-twitch
ExecStart=/home/snaxine/.local/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ika-twitch

[Install]
WantedBy=multi-user.target
